# services/sendspin-watchdog.service
[Unit]
Description=Sendspin stale stream watchdog
After=sound.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sendspin-watchdog.sh
