# services/sendspin@.service
# PRIMARY: Sendspin connects to lox-audioserver for Loxone multiroom audio
# Each room gets its own service instance: sendspin@room_kitchen, etc.
# The --audio-device must match an ALSA PCM device from generate_alsa_config.py
[Unit]
Description=Sendspin client for %i
After=sound.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/tobias/.local/bin/sendspin daemon --url ws://localhost:7090/sendspin --id %i --name %i --audio-device %i
Restart=on-failure
RestartSec=5
User=tobias

[Install]
WantedBy=multi-user.target
