# services/snapclient@.service
# LEGACY: Snapcast setup - only needed if not using lox-audioserver
# Port 1714 is used to avoid conflict with lox-audioserver's built-in port 1704
[Unit]
Description=Snapcast client for %i
After=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/snapclient -h localhost -p 1714 --hostID %i -s %i --logsink=system --sampleformat 48000:16:*
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
